# root/etc/s6-overlay/s6-rc.d/run-mod-code-server-traefik/up
/etc/s6-overlay/s6-rc.d/run-mod-code-server-traefik/run